# .config/autostart/Devilspie.desktop
[Desktop Entry]
Type=Application
Name=Devilspie
Description=
Exec=devilspie2
