# .config/autostart/Invidious.desktop
[Desktop Entry]
Type=Application
Name=Invidious
Description=
Exec=bash -c invidious
