# .config/autostart/Podsync.desktop
[Desktop Entry]
Type=Application
Name=Podsync
Description=
Exec=bash -c podsync
