# autostart/syncthing.desktop
[Desktop Entry]
Type=Application
Name=syncthing
Description=
Exec=syncthing
