# autostart/syncthing.desktop
[Desktop Entry]
Exec=syncthing
Type=Application
