# autostart/xremap.desktop
[Desktop Entry]
Type=Application
Name=xremap
Description=
Exec=fish -c remap
