# autostart/xremap.desktop
[Desktop Entry]
Type=Application
Name=xremap
Description=
Exec=xremap --mouse /home/glados/.config/xremap.yml
