# autostart/xremap.desktop
[Desktop Entry]
Exec=fish -c remap
Type=Application
