# dot-config/autostart/invidious.desktop
[Desktop Entry]
Type=Application
Name=Invidious
Description=
Exec=fish -c invid
